# ARCHLOGS/bin/archmonitord
#!/bin/bash
# chkconfig: 2345 90 15
# description: 
# daemon to run /usr/local/bin/archmonitor.sh
# which is i symbolic link to /opt/netapp/snapcreator/scripts/ARCHLOGS/bin/archmonitor.sh

# Source function library.
. /etc/rc.d/init.d/functions

### installation path
mypath=/opt/netapp/snapcreator/scripts/ARCHLOGS

### END installation path
. $mypath/etc/archivelogs_backup.conf || exit 1


test -f /usr/local/bin/archmonitor.sh || (echo "cannot found /usr/local/bin/archmonitor.sh"; kill $$)

start() {
    if [ ! -d "$path_to_archlogs_dir" ]
    then
     echo "ERROR: ARCHLOGS directory $path_to_archlogs_dir not found!"
     exit 1
    fi
    getstatus=`/etc/init.d/archmonitord status | grep Running | awk '{print $1}' `
    if [ "$getstatus" == "Running" ]
    then
     echo "archmonitord is already running"
     exit 1
    else
     echo "Starting archmonitord:"
     /usr/local/bin/archmonitor.sh &
     sleep 5
     /etc/init.d/archmonitord status
    fi
}

stop() {
     echo "Shutting down archmonitord:"
     sleep 1
     proc=`pgrep -f archmonitor.sh`
     if test -z "$proc"; then
      echo "archmonitord failed to stop, not running!"
      exit 1
     else
      kill -9 $proc
     fi
     result=$?
     if [ $result == 0 ]
     then
      echo "archmonitord stopped"
     else
      echo "archmonitord failed to stop!"
      exit 1
     fi
     
}

case "$1" in 
    start)
       start
       ;;
    stop)
       stop
       ;;
    restart)
       stop
       start
       ;;
    status)
       echo "Checking Status of archmonitord: "
       sleep 1
       proc=`pgrep -f archmonitor.sh`
       if test -z "$proc"; then
        echo "Stopped"
       else
        echo "Running"
       fi
       ;;
    jobreview)
       getstatus=`/etc/init.d/archmonitord status | grep Running | awk '{print $1}'`
       if [ "$getstatus" != "Running" ]
       then
        echo "archmonitord isn't running"
        exit 1
       else
        /usr/local/bin/archmonitor.sh --jobreview
       fi
       ;;
    *)
       echo "Usage: $0 {start|stop|status|restart|jobreview}"
esac

exit 0 
